# example/110_gentoo/openrc/config/user-pipewire.initd
#!/sbin/openrc-run

description="PipeWire media server"
command="/usr/bin/pipewire"
command_background="yes"
pidfile="${XDG_RUNTIME_DIR}/pipewire.pid"

depend() {
  need dbus
}
