# example/110_gentoo/docker_min_hpz6_openrc/config/reverse-ssh-us.initd
#!/sbin/openrc-run

name="reverse-ssh-us"
description="Reverse SSH tunnel to tinyus"
supervisor="supervise-daemon"
command="/usr/bin/autossh"
command_args="-M 0 -N -o ExitOnForwardFailure=yes -o ServerAliveInterval=30 -o ServerAliveCountMax=3 -R 2332:localhost:22 tinyus"
command_user="kiel"
supervise_daemon_args="-e HOME=/home/kiel"

depend() {
    use net
    after iwd
    want tailscale
}
