# example/110_gentoo/docker_min_hpz6_openrc/config/user-pipewire-pulse.initd
#!/sbin/openrc-run

description="PipeWire PulseAudio compatibility server"
command="/usr/bin/pipewire-pulse"
command_background="yes"
pidfile="${XDG_RUNTIME_DIR}/pipewire-pulse.pid"

depend() {
  need pipewire dbus
}
