# example/110_gentoo/docker_min_hpz6_openrc/config/user-wireplumber.initd
#!/sbin/openrc-run

description="WirePlumber session manager"
command="/usr/bin/wireplumber"
command_background="yes"
pidfile="${XDG_RUNTIME_DIR}/wireplumber.pid"

depend() {
  need pipewire dbus
}
